# command/hourly_up_snippet.crontab
0 * * * *  cd ~/agideo-snippets && git pull > /dev/null 2>&1
